# opsell-direction.service
[Unit]
Description=Market Direction Classifier Microservice
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/opsell
Environment="SERVICE_HOST=0.0.0.0"
Environment="SERVICE_PORT=8000"
Environment="SERVICE_WORKERS=4"
Environment="PYTHONPATH=/opt/opsell"

# Python virtual environment
Environment="PATH=/opt/opsell/venv/bin:$PATH"
ExecStart=/opt/opsell/venv/bin/gunicorn \
    --chdir /opt/opsell \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 0.0.0.0:8000 \
    --timeout 120 \
    --access-logfile /var/log/opsell/access.log \
    --error-logfile /var/log/opsell/error.log \
    --log-level info \
    services.direction:app

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
